# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:lagrange]
command=/app/Lagrange.OneBot
directory=/app
autostart=true
autorestart=true
startretries=3
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DOTNET_SYSTEM_GLOBALIZATION_INVARIANT="1"

[program:astrbot]
command=python /app/main.py
directory=/app
autostart=true
autorestart=true
startretries=3
priority=2
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/root/.pyffmpeg/bin"

# 进程组配置，方便统一管理
[group:bot]
programs=lagrange,astrbot
priority=999

# 示例：如果你想添加更多程序，可以像下面这样配置
# [program:example_service]
# command=/path/to/your/program
# directory=/path/to/working/directory
# autostart=true
# autorestart=true
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0

# 示例：Node.js应用
# [program:nodeapp]
# command=node /app/server.js
# directory=/app
# autostart=true
# autorestart=true
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0

# 示例：后台任务
# [program:worker]
# command=python /app/worker.py
# directory=/app
# autostart=true
# autorestart=true
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
